# rootfs/etc/services.d/knx-failover/finish
#!/bin/bash
# This script runs when the service is being stopped
exit 0
