# rootfs/etc/services.d/knx-haproxy/run
#!/command/execlineb -P
define LOG_LEVEL info
/run.sh
